# packages/dci-umb/dci-umb-0.4.2.post202510211220.tar.gz/dci-umb-0.4.2.post202510211220/systemd/dci-umb.service
[Service]
Type=simple
EnvironmentFile=-/etc/dci-umb/config
ExecStart=/usr/bin/dci-umb
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
