# packages/meltingplot.rpi-camera/meltingplot_rpi_camera-0.0.6-py3-none-any.whl/meltingplot_rpi_camera-0.0.6.data/data/rpi-camera.service
[Unit]
Description=RPi Camera MJPEG Server
After=network.target

[Service]
Type=simple
Restart=always
StartLimitInterval=0
StartLimitBurst=1440
RestartSec=60
ExecStart=/usr/local/bin/rpi-camera start
WorkingDirectory=/home/pi
User=pi
Group=pi
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target