# packages/sleepless-agent/sleepless_agent-0.1.0.tar.gz/sleepless_agent-0.1.0/com.sleepless-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.sleepless-agent</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3</string>
		<string>-m</string>
		<string>sleepless_agent.daemon</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/YOUR_USERNAME/projects/sleepless-agent</string>
	<key>StandardOutPath</key>
	<string>/Users/YOUR_USERNAME/projects/sleepless-agent/logs/daemon.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOUR_USERNAME/projects/sleepless-agent/logs/daemon.err</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Restart</key>
	<integer>10</integer>
</dict>
</plist>
